# debian/postinst
#!/bin/bash
#
# postinst script for linuxmuster-relution
# thomas@linuxmuster.net
# 20180312
# GPL v3
#

# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

name="linuxmuster-relution"
ngxname="linuxmuster-nginx-proxy"
images="relution/relution:latest"
sharedir="/usr/share/$name"
dockerkey="/etc/linuxmuster/ssl/docker.key.pem"
dockerdir="/srv/docker/$name"
cmpstpl="$sharedir/docker-compose.yml"
cmpsyml="$dockerdir/docker-compose.yml"
rltntpl="$sharedir/$name.conf"
rltncfg="/srv/docker/$ngxname/conf.d/$name.conf"
interfaces="/etc/network/interfaces"

case "$1" in

 configure)

  # pull docker image
  if ping -q -c2 google.com &> /dev/null; then
    for i in $images; do
      docker pull "$i"
    done
  fi

  if [ ! -e "$dockerkey" ]; then
    echo "SSL hostkey not found!"
    echo "Please run linuxmuster-prepare with -c option first,"
    echo "then run dpkg-reconfigure $name."
    exit 0
  fi

  # get hostip and hostname
  hostip="$(grep address "$interfaces" | head -1 | awk '{ print $2 }')"
  fqdn="$(cat /etc/hostname)"
  hname="$(echo "$fqdn" | awk -F. '{ print $1 }')"
  dname="${fqdn/$hname./}"

  # create configuration
  sed -e "s|@@hostname@@|$hname|
          s|@@domainname@@|$dname|" "$cmpstpl" > "$cmpsyml"
  sed -e "s|@@hostip@@|$hostip|" "$rltntpl" > "$rltncfg"

  # enable and start service
  /bin/systemctl enable "$name.service"
  /bin/systemctl stop "$name.service"
  /bin/systemctl start "$name.service"
  /bin/systemctl restart "$ngxname.service"

 ;;

 abort-upgrade|abort-remove|abort-deconfigure)
  rm -f "$rltncfg"
 ;;

 *)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
 ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
